# conf/supervisord.conf.d/beat.conf
; ================================
;  celery beat supervisor example
; ================================

[program:celerybeat]
; Set full path to celery program if using virtualenv
command=celery beat -A webproject --pidfile=/opt/publisher/run/beat.pid --loglevel=INFO --schedule=/opt/publisher/run/beat.db

; remove the -A myapp argument if you are not using an app instance

stdout_logfile=/var/log/publisher/beat.log
stderr_logfile=/var/log/publisher/beat.log
directory=/opt/publisher
numprocs=1
autostart=true
autorestart=true
startsecs=10
user=publisher

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=999
